# upgrade.service
[Unit]
Description=Execute a system upgrade.

[Service]
Type=simple
ExecStart=/usr/local/bin/upgrade
